# sentinel.service
[Unit]
Description=Sentinel — AI Agency Orchestrator
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/sentinel
ExecStart=/usr/bin/python3 /opt/sentinel/sentinel.py scheduler
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
